# supervisor-example.conf
[program:battle-mages]
command=/usr/bin/php {PATH_TO_PROJECT} artisan queue:work
numprocs=1
autostart=true
autorestart=true
user=user
stdout_logfile={PATH_TO_PROJECT}/storage/logs/supervisor.txt
